# package/server/deb/src/deb/control/postinst
#!/bin/sh

systemctl daemon-reload
echo "Starting Intuso Housemate"
systemctl enable housemate
systemctl start housemate

exit 0