# package/server/deb/src/deb/control/postinst
#!/bin/sh

systemctl daemon-reload
echo "Starting Intuso Housemate Server"
systemctl enable housemate-server
systemctl start housemate-server

exit 0