# package/node-1-0/deb/src/deb/control/postinst
#!/bin/sh

systemctl daemon-reload
echo "Starting Intuso Housemate Node"
systemctl enable housemate-node
systemctl start housemate-node

exit 0